# deploy/systemd/profeta-api.service
[Unit]
Description=Profeta API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=profeta
Group=profeta
WorkingDirectory=/opt/profeta/current
EnvironmentFile=/opt/profeta/current/.env
Environment=PYTHONUNBUFFERED=1
RuntimeDirectory=profeta
ExecStart=/bin/bash /opt/profeta/current/scripts/start_api.sh
Restart=always
RestartSec=5
TimeoutStartSec=45
TimeoutStopSec=30
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
